# DEV05-SAM9260/ROOTFS_TREE/etc/.sysconfig/crontabs/root
SHELL=/bin/sh
PATH=/sbin:/bin:/usr/sbin:/usr/bin
HOME=/

# run-parts
0 0 * * * /bin/sh /etc/sysconfig/cronjobs/ntpd 1
0 */12 * * * /bin/sh /etc/sysconfig/cronjobs/synclock

0 0 * * * /etc/init.d/S06_ckver

