# DEV05-SAM9260/ROOTFS_TREE/etc/inittab
# /etc/inittab
#
# Copyright (C) 2001 Erik Andersen <andersen@codepoet.org>
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console. 
#              If specified, then /dev/$id device must exist
# runlevels == ignored, busybox doesn't support it
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# Startup the system
# mount all the file systems specified in /etc/fstab 
::sysinit:/bin/mount -a

#Use mdev as hotplug to auto mount USB storage or SD card 
::sysinit:/bin/echo /sbin/mdev > /proc/sys/kernel/hotplug
::sysinit:/sbin/mdev -s

#make shm, pts support
::sysinit:/bin/mkdir -p /dev/pts
::sysinit:/bin/mkdir -p /dev/shm
::sysinit:/bin/mount -t devpts devpts /dev/pts

#Mount our apps/info partition
null::wait:/bin/mount -o sync,noatime,ro -t jffs2 /dev/mtdblock6 /apps
null::wait:/bin/mount -o sync,noatime,ro -t jffs2 /dev/mtdblock7 /info

null::once:/bin/mkdir -p /tmp/logs
null::once:/bin/mkdir -p /tmp/stat
null::once:/bin/touch /tmp/lastlog
null::once:/bin/touch /tmp/wtmp
#null::wait:/sbin/ubiattach /dev/ubi_ctrl -m9 -d 9

# now run any rc scripts
null::wait:/etc/init.d/rcS

#Set hostname 
null::once:/bin/hostname -F /etc/hostname
null::once:/usr/sbin/powerTime

#Enable console logon
null::respawn:/sbin/getty -L ttyS0 115200 vt100

# system daemo
null::once:/sbin/syslogd -n
#null::once:/sbin/klogd -n

# Stuff to do before rebooting
null::shutdown:/usr/sbin/lastTime
null::shutdown:/usr/sbin/powerTime END
null::shutdown:/bin/umount /apps
null::shutdown:/bin/umount /info
null::shutdown:/bin/umount /data
#null::shutdown:/bin/killall klogd
null::shutdown:/bin/killall syslogd
null::shutdown:/bin/umount -a -r

